# docker/platform/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi.ini --enable-threads

[program:nginx-app]
command = /usr/sbin/nginx